# config/production/supervisor.conf
[program:pdepou]
command = /home/ubuntu/pdepou/config/production/gunicorn.sh
user = root
stdout_logfile = /var/log/supervisor/pdepou.log
redirect_stderr = true
